# distributions/deb-package/snoozenode/debian/postrm
#!/bin/bash

SCRIPT_NAME="snoozenode"
update-rc.d $SCRIPT_NAME remove
rm -Rf $SCRIPT_NAME
